# roles/recovery_slack_listener/templates/recovery-slack-listener.service.j2
[Unit]
Description=Homelab Recovery Slack Listener
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ recovery_slack_listener_user }}
EnvironmentFile={{ recovery_slack_listener_key_dir }}/slack-listener.env
ExecStart={{ recovery_slack_listener_install_dir }}/venv/bin/python {{ recovery_slack_listener_script }}
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
