# roles/systemd_timers/templates/ansible-timer.timer.j2
[Unit]
Description={{ timer.description }} (timer)

[Timer]
OnCalendar={{ timer.schedule }}
Persistent={{ timer.persistent | ternary('true', 'false') }}

[Install]
WantedBy=timers.target
