# roles/recovery_io/templates/recovery-io.service.j2
[Unit]
Description=Homelab Recovery I/O (Slack listener)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ recovery_io_user }}
Group={{ recovery_io_user }}
ExecStart={{ recovery_io_install_dir }}/venv/bin/python {{ recovery_io_script }}
EnvironmentFile={{ recovery_io_env_dir }}/recovery-io.env
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Hardening (§7)
# NoNewPrivileges left unset (default=no) — recovery-io must call sudo.
PrivateTmp=yes
ProtectSystem=strict
# read-only + ReadWritePaths: Codex (child process via sudo) needs to write
# to ~/.codex/ for its in-process app-server state. DAC (0700 recovery-exec)
# still prevents recovery-io user from accessing those files.
ProtectHome=read-only
ReadWritePaths={{ recovery_io_exec_home }}
RestrictSUIDSGID=yes
LockPersonality=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

[Install]
WantedBy=multi-user.target
